# desk/docker/dns/service/pdns/run
#!/bin/sh
exec /usr/sbin/pdns_server --daemon=no --guardian=no
# >> /var/log/pdns.log 2>&1
