# desk/images/dns/etc/service/worker/run
#!/bin/sh
cd /projects/desk/desk
exec /var/py27/bin/python ./dworker run -v
