# desk/docker/worker/files/service/worker/run
#!/bin/sh
cd /opt/app/desk
exec /var/py27/bin/python -u ./dworker run -v >> -WORKER_LOG-/worker.log 2>&1
