# desk/images/foreman/etc/service/worker/run
#!/bin/sh
cd /projects/desk/desk && /var/py27/bin/python ./dworker run
