# desk/docker/worker/etc/services.d/worker/run
#!/usr/bin/with-contenv sh
#set -e

cd /opt/app/desk
echo "* staring worker"
sleep 2 # start after other /etc/services.d
python -u ./dworker run -v >> "${WORKER_LOG}/${FQHOSTNAME}.log" 2>&1
