# desk/images/couchdb/etc/service/couchdb/run
#!/bin/sh
export HOME=/usr/local/var/lib/couchdb/
exec /usr/local/bin/couchdb
