# desk/docker/dns/files/service/pdns/run
#!/bin/sh
exec /usr/sbin/pdns_server --daemon=no --guardian=no >> -PDNS_LOG-/pdns.log 2>&1
